# root/etc/services.d/mysql/run
#!/bin/sh
exec /usr/bin/mysqld --user=mysql --console
